# sharemechat-v1/ops/prod-access-normalizer/systemd/sharemechat-prod-access-normalizer.timer
[Unit]
Description=Run PROD access log normalization every 5 minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
Unit=sharemechat-prod-access-normalizer.service

[Install]
WantedBy=timers.target
