# sharemechat-v1/ops/audit-access-normalizer/systemd/sharemechat-audit-access-normalizer.timer
[Unit]
Description=Run AUDIT access log normalization every 5 minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
Unit=sharemechat-audit-access-normalizer.service

[Install]
WantedBy=timers.target
