# sharemechat-v1/ops/audit-access-normalizer/systemd/sharemechat-audit-access-normalizer.service
[Unit]
Description=Normalize AUDIT access logs from CloudFront and Nginx
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/sharemechat-audit-access-normalizer/bin/normalize-audit-access.sh --config /etc/sharemechat-audit-access-normalizer/config.env
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
