# sharemechat-v1/ops/test-access-blocker/systemd/sharemechat-test-access-blocker.service
[Unit]
Description=Propose deny list for TEST based on classifier summary (DRY-RUN)
# El classifier de TEST no es una unit independiente: vive como sub-paso
# dentro de sharemechat-test-daily-report.service (07:10 UTC). Apuntamos
# After= a la unit real para que la dependencia ordene correctamente cuando
# ambas estan encoladas. Patron corregido en paquete 10.B.3 (la referencia
# anterior a sharemechat-test-access-classifier.service apuntaba a una unit
# que nunca existio; systemd ignoraba el After= silenciosamente).
# Nota: aun corregida la dependencia, el race condition con catch-up agresivo
# tras rearrancar TEST queda mitigado por el skip elegante del blocker:
# si el summary del dia esperado no existe, exit 0 con log informativo.
After=network-online.target sharemechat-test-daily-report.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
# DRY-RUN ONLY. This unit does NOT reload nginx and does NOT write
# /etc/nginx/deny-test-ips.conf. It only produces advisory artifacts
# under /var/log/sharemechat-test-access-blocker/.
ExecStart=/opt/sharemechat-test-access-blocker/bin/block-test-access.sh --config /etc/sharemechat-test-access-blocker/config.env
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
